# MySQL/crontab
@reboot ~/Desktop/WateringSystem-server/MySQL/mqtt_mysql_handler.sh
*/15 * * * * ~/Desktop/WateringSystem-server/MySQL/mqtt_mysql_handler.sh
